# deploy/tk-selection-backend.service
[Unit]
Description=TK Selection Backend
After=network.target

[Service]
Type=simple
User=tkbackend
Group=tkbackend
WorkingDirectory=/opt/tk-selection/backend
EnvironmentFile=/etc/tk-selection/backend.env
Environment=PLAYWRIGHT_BROWSERS_PATH=/opt/ms-playwright
ExecStart=/opt/tk-selection/backend/.venv/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 1
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
